# playbooks/roles/wordpress/templates/rc.local.j2
#!/bin/bash
cd /vagrant/cakephp 
sudo -u vagrant grunt watch &